# supervisor/jenkins.conf
[program:jenkins]
autorestart = true
command = /usr/local/bin/jenkins.sh
priority = 2
user = jenkins
