# supervisor/jenkins.conf
[program:jenkins]
autorestart = true
command = /usr/local/bin/jenkins.sh
environment=HOME=%(ENV_JENKINS_HOME)s
priority = 2
user = jenkins
